# conf/supervisor.d/php-fpm-exporter.ini
[program:php-fpm-exporter]
command = /usr/local/bin/php-fpm-exporter server --phpfpm.scrape-uri "unix:///dev/shm/php-cgi.sock;/fpm-status" --phpfpm.fix-process-count true
autostart=false
autorestart=false
priority=992
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
